# layout/devops/deploy/sphere.service
[Unit]
Description=Sphere Application
After=network.target

[Service]
Type=simple
ExecStart=/opt/sphere/app --config /opt/sphere/config.json
WorkingDirectory=/opt/sphere
User=nobody
Group=nobody
Restart=always
RestartSec=5
Environment="APP_ENV=production"
EnvironmentFile=/opt/sphere/.env
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target